# systemctl/freeark-plc-connection-monitor.service
[Unit]
Description=FreeArk PLC Connection Monitor Service
Documentation=https://github.com/your-repo/freeark
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=5

# 定义环境变量
Environment=PYTHONPATH=/home/yangyang/Freeark/FreeArk

# 工作目录
WorkingDirectory=/home/yangyang/Freeark/FreeArk/

# 启动命令 - 使用虚拟环境中的Python
# 检查间隔: 300秒(5分钟)
# 超时阈值: 600秒(10分钟)
ExecStart=/home/yangyang/Freeark/FreeArk/venv/bin/python /home/yangyang/Freeark/FreeArk/FreeArkWeb/backend/freearkweb/manage.py plc_connection_monitor --check-interval 300 --timeout-threshold 600

# 停止命令
ExecStop=/bin/kill -s SIGINT $MAINPID

# 设置用户和组（根据实际部署需求调整）
User=yangyang

# 日志输出
StandardOutput=journal
StandardError=journal
SyslogIdentifier=freeark-plc-connection-monitor

[Install]
WantedBy=multi-user.target