# systemctl/freeark-backend.service
[Unit]
Description=FreeArk backend Service
After=network.target redis-server.service

[Service]
Type=simple
User=yangyang
WorkingDirectory=/home/yangyang/Freeark/FreeArk/
# 使 Python stdout/stderr 不缓冲，启动日志（含 Uvicorn 启动参数）实时写入 journal
Environment=PYTHONUNBUFFERED=1
# Uvicorn ASGI（替换 Waitress，2026-05-23 OpenClaw 集成 + WebSocket 支持，见 ADR-001）
# perf-P1a（2026-05-31）：--workers 1 → --workers 2
#   原因：换用 RedisChannelLayer（channels_redis），解除 InMemoryChannelLayer 单进程限制。
#   Pi 5 aarch64 4 核，2 worker 约占 200MB，留足余量给 OpenClaw/MySQL/OS。
#   依赖：redis-server（apt 安装，loopback 127.0.0.1:6379，开机自启）。
# 回滚到单 worker：将 --workers 2 改回 --workers 1，同时将 settings.py CHANNEL_LAYERS
#   的 BACKEND 改回 "channels.layers.InMemoryChannelLayer"（见 ADR-P1A-001 回滚方案）。
# 回滚到 Waitress: 将下行 ExecStart 改回 `python <repo>/FreeArkWeb/backend/freearkweb/start_waitress_server.py`
ExecStart=/home/yangyang/Freeark/FreeArk/venv/bin/uvicorn freearkweb.asgi:application --host 0.0.0.0 --port 8000 --workers 2 --app-dir /home/yangyang/Freeark/FreeArk/FreeArkWeb/backend/freearkweb
Restart=on-failure
RestartSec=30s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=freeark-backend

[Install]
WantedBy=multi-user.target