# deployment/systemd/freeark-condensation-consumer.service
# freeark-condensation-consumer.service — 结露预警 MQTT 消费者（v0.7.0-CW，MOD-INFRA-CW-01）
#
# 部署步骤（在生产服务器 192.168.31.51 执行）：
#   sudo cp deployment/systemd/freeark-condensation-consumer.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now freeark-condensation-consumer
#
# 验证：
#   journalctl -u freeark-condensation-consumer -n 30 -f
#   systemctl status freeark-condensation-consumer

[Unit]
Description=FreeArk Condensation Warning MQTT Consumer (v0.7.0-CW)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=yangyang
WorkingDirectory=/home/yangyang/Freeark/FreeArk/FreeArkWeb/backend/freearkweb
ExecStart=/home/yangyang/Freeark/FreeArk/venv/bin/python manage.py condensation_consumer
Restart=on-failure
RestartSec=30s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=freeark-condensation-consumer
Environment=DJANGO_SETTINGS_MODULE=freearkweb.settings

[Install]
WantedBy=multi-user.target
